# systemd-services/beacon-chain.service
[Unit]
Description     = eth2 beacon chain service
Wants           = network-online.target
After           = network-online.target

[Service]
User            = ethereum
WorkingDirectory= /home/ethereum/lighthouse
EnvironmentFile = /home/ethereum/testnet/nodevars_env.txt
ExecStart       = /home/ethereum/lighthouse/lighthouse bn --testnet-dir /home/ethereum/testnet --datadir /home/ethereum/data-lh --staking --metrics --port 9010 --http-address 0.0.0.0 --metrics-address 0.0.0.0 --execution-endpoint http://127.0.0.1:8551 --execution-jwt="/home/ethereum/jwtsecret" --boot-nodes ${BOOTNODE_ENR_LIST}
Restart         = on-failure

[Install]
WantedBy= multi-user.target